# install-files/profile.d/libreoffice.sh

# use GTK2 look and feel, instead of GTK3.
# Especially useful on 18.04 where libreoffice-gtk3 is a dependency
# of libreoffice-gnomei, which is needed for samba shares.
export SAL_USE_VCLPLUGIN=gtk
